# tvheadend.conf
[program:tvheadend]
user = nobody
command = /usr/bin/tvheadend -f -p /run/tvheadend.pid -C -u hts -g video