# tvheadend.conf
[program:tvheadend]
autorestart = false
startsecs = 0
user = %(ENV_USERNAME)s
command = "/usr/bin/tvheadend" -f -p "/run/tvheadend.pid" -C -u %(ENV_USERNAME)s -g %(ENV_GROUPNAME)s