# tvheadend.conf
[program:tvheadend]
user = %(ENV_UID)s
command = "/usr/bin/tvheadend" -f -p "/run/tvheadend.pid" -C -u %(ENV_UID)s -g %(ENV_GID)s